# code/ansible/playbooks/files/PruneDocker.service
[Unit]
Description=PruneDocker
[Service]
Type=oneshot
ExecStart=/usr/bin/docker image prune -f
WorkingDirectory=/tmp